# powergslb/powergslb.service
# PowerGSLB systemd service

[Unit]
Description=PowerGSLB service

[Service]
Type=notify
ExecStart=/usr/sbin/powergslb -c /etc/powergslb/powergslb.conf
WatchdogSec=3
Restart=on-failure
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_RAW
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
