# docker/rootfs/usr/lib/systemd/system/mariadb.service.d/mariadb-install-db.conf
[Service]
ExecStartPre=+/bin/sh -c 'test -d /var/lib/mysql/mysql || \
    mariadb-install-db --datadir=/var/lib/mysql --user=mysql --group=mysql --skip-test-db'
